# browser-bot.service
[Unit]
Description=Browser Bot Service
After=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/opt/browser_bot
Environment="NODE_ENV=production"
EnvironmentFile=/opt/browser_bot/.env
ExecStart=/usr/bin/node /opt/browser_bot/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=browser-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/browser_bot/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=4G

[Install]
WantedBy=multi-user.target

